# bin/healthcheck.cron.template
* * * * * <F1SIM_HOME>/bin/healthcheck.sh
